# RootHideSMSImport/layout/Library/LaunchDaemons/com.nextsolution.dailyledger.smsimport.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nextsolution.dailyledger.smsimport</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/libexec/dailyledgersmsd</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ProcessType</key>
    <string>Background</string>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>UserName</key>
    <string>mobile</string>
</dict>
</plist>
